# dist/.deb-build/uxplay-tray_3cc8bf6_amd64/etc/xdg/autostart/uxplay-tray.desktop
[Desktop Entry]
Type=Application
Name=UxPlay Tray Background Service
Comment=Start UxPlay Tray in the background when you log in
Exec=/usr/bin/uxplay-tray
Icon=uxplay-tray
Terminal=false
NoDisplay=true
OnlyShowIn=GNOME;Unity;
X-GNOME-Autostart-enabled=true
X-GNOME-Autostart-Delay=3
